# app/automation-runner.service
[Service]
Type=simple
WorkingDirectory=/opt/automation-runner

EnvironmentFile=/opt/automation-runner/.env
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/opt/automation-runner
Environment=XDG_CONFIG_HOME=/opt/automation-runner/.config

ExecStart=/opt/automation-runner/.venv/bin/python /opt/automation-runner/app/app.py

Restart=on-failure
RestartSec=5